# external/config/autostart/yakuake.desktop
[Desktop Entry]
Type=Application
Exec=yakuake
Icon=yakuake
Name=Yakuake
Comment=Quake-style dropdown terminal hosting nix
OnlyShowIn=KDE;
X-KDE-autostart-after=panel
X-GNOME-Autostart-enabled=true
